# etc/flydigid.service
[Unit]
Description=Daemon for flydigictl

[Service]
BusName=com.pipe01.flydigi.Gamepad
ExecStart=/usr/bin/flydigid
